# root/etc/cron.hourly/httpd-php-revproxy_creds.cron
#!/bin/bash

curl -k -s 'https://'`uname -n`'/api/cred_cron.php'
